# LaunchAgents/com.example.yapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template launchd agent for the Yapper TTS server.

  Copy to ~/Library/LaunchAgents/, edit the four __REPLACE__ values for your
  home directory, then:

    launchctl bootstrap "gui/$(id -u)" ~/Library/LaunchAgents/com.example.yapper.plist
    launchctl kickstart -k "gui/$(id -u)/com.example.yapper"

  Logs: ~/Library/Logs/yapper/
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.yapper</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPLACE__/.local/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>__REPLACE__/Documents/yapper</string>
        <string>python</string>
        <string>server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPLACE__/Documents/yapper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>YAPPER_PORT</key>
        <string>8765</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__REPLACE__/Library/Logs/yapper/server.out.log</string>

    <key>StandardErrorPath</key>
    <string>__REPLACE__/Library/Logs/yapper/server.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
